# apps/lumiere-update/lumiere-update.service
[Unit]
Description=Lumière OS Güncelleme Kontrolü Servisi

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lumiere-update check
User=root
